# deployment/host_archive/etc/supervisor/conf.d/pb-blue.conf
[program:pb-blue]
command=/opt/peerbook/blue/peerbook --addr 127.0.0.1:40000 
environment=
	PB_HOME_URL="http://%(ENV_PUBLIC_IP)s:8000",
	PB_ICETCP_PORT="40001"
